# supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock   ; Chemin vers le socket

[supervisord]
logfile=/var/log/supervisord.log ; journal d'activité
pidfile=/var/run/supervisord.pid ; PID file
childlogdir=/var/log/supervisor  ; dossiers pour les logs

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:bot]
command=/root/entrypoint.sh
autostart=true
autorestart=true
